# one-united-ts.service
[Unit]
Description=one-united-ts - LLM API Gateway
Documentation=https://github.com/JmPotato/one-united-ts
After=network.target

[Service]
WorkingDirectory=/opt/one-united-ts
ExecStart=bun src/app.ts
Restart=on-failure
RestartSec=5

DynamicUser=yes
StateDirectory=one-united-ts

Environment=NODE_ENV=production
Environment=PORT=5222
Environment=HOSTNAME=127.0.0.1
Environment=DATABASE_PATH=/var/lib/one-united-ts/kv.db
Environment=ONE_API_KEY=sk-your-key-here

StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=true
DeviceAllow=/dev/null rw

[Install]
WantedBy=multi-user.target